# ninjakeys.conf
description 	"NinjaKeys"
author		"Mario Pareja <pareja.mario@gmail.com>"

start on started

script
	cd /opt/ninjakeys
	exec ./run.sh
end script
