# touchcursor.conf
description "TouchCursor"
author		"Mario Pareja <pareja.mario@gmail.com>"

start on started

script
	cd /opt/kbd-mangler
	exec ./run.sh tc.js
end script
